# deploy/systemd/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=HOME=%h
ExecStart=%h/agri-orchestrator-main/scripts/run-openclaw-gateway.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
